# defaultconfig/rhasspy/profiles/en_base/supervisord.conf
[supervisord]
nodaemon=true

[program:mqtt]
command=mosquitto -p 12183
priority=0
stopasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:speech_to_text]
command=rhasspy-asr-kaldi-hermes --model-type nnet3 --model-dir /profiles/en/kaldi/model --graph-dir /profiles/en/kaldi/model/graph --debug --host localhost --port 12183 --site-id speaker --site-id speaker --dictionary /profiles/en/kaldi/dictionary.txt --language-model /profiles/en/kaldi/language_model.txt --language-model-type text_fst --base-dictionary /profiles/en/kaldi/base_dictionary.txt --base-dictionary /profiles/en/kaldi_custom_words.txt --dictionary-casing lower --g2p-model /profiles/en/kaldi/g2p.fst --unknown-words /profiles/en/kaldi_unknown_words.txt --base-language-model-fst /profiles/en/kaldi/base_language_model.fst --base-language-model-weight 0 --mixed-language-model-fst /profiles/en/kaldi/mixed.fst --voice-skip-seconds 0 --voice-min-seconds 1 --voice-max-seconds 20 --voice-speech-seconds 0.3 --voice-silence-seconds 0.5 --voice-before-seconds 0.5 --voice-sensitivity 1 --voice-silence-method vad_only
stopasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:intent_recognition]
command=rhasspy-nlu-hermes --intent-graph /profiles/en/intent_graph.pickle.gz --debug --host localhost --port 12183 --site-id speaker --site-id speaker --replace-numbers --language en_US --casing lower --converters-dir /profiles/en/converters
stopasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:intent_handling]
command=rhasspy-homeassistant-hermes --url http://192.168.0.159:8123/ --debug --host localhost --port 12183 --site-id speaker --site-id speaker --handle-type intent --access-token eyJ0eXAiOiJKV1QiLCJhbGciOiJIUzI1NiJ9.eyJpc3MiOiI0NWM5NTAwMjlkYjg0MjZmOTZlYTE0M2ZlMjQ0ODFjYyIsImlhdCI6MTYyNDAwMzYwMiwiZXhwIjoxOTM5MzYzNjAyfQ.mi8QwUF5EXCGaWUMynZ4Rn7RLMv2RazTa44-NGaf16c --event-type-format rhasspy_{0}
stopasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:text_to_speech]
command=rhasspy-tts-cli-hermes --tts-command 'espeak --stdout -v {lang}' --voices-command 'espeak --voices | tail -n +2 | awk '"'"'{ print $2,$4 }'"'"'' --language en-gb --volume 1 --debug --host localhost --port 12183 --site-id speaker
stopasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

